# env/etc/udev/rules.d/fstab.rules
# Force fstab options for devices
RUN+="/usr/sbin/rbfstab"
